# packaging/systemd/debug/kismet-debug.service.in
[Unit]
Description=Kismet debug
ConditionPathExists=@prefix@/bin/kismet_debug

[Service]
User=root
Group=root
Type=simple
ExecStart=@prefix@/bin/kismet_debug
TimeoutSec=0
SysVStartPriority=99
Restart=always

[Install]
WantedBy=multi-user.target

