# packaging/systemd/kismet.service.in
[Unit]
Description=Kismet
ConditionPathExists=@prefix@/bin/kismet
After=network.target auditd.service

[Service]
User=root
Group=root
Type=simple
ExecStart=@prefix@/bin/kismet
KillMode=process
TimeoutSec=0
SysVStartPriority=99
Restart=always

[Install]
WantedBy=multi-user.target

